# infra/ansible/roles/sharingan-node/templates/madara-full-node.service.j2
[Unit]
Description=Madara Full Node

[Service]
WorkingDirectory=/home/admin/madara
Restart=on-failure
ExecStart=/home/admin/madara/target/release/madara --testnet sharingan

ExecReload=/bin/kill -HUP
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
